# bbs.ini
[program:mbtask]
command=%(ENV_MBSE_ROOT)s/bin/mbtask -nd
user=mbse

[program:xinetd]
command=/usr/sbin/xinetd -dontfork

[program:darkhttpd]
command=/usr/sbin/darkhttpd /srv/mbse/www/htdocs
