# xinetd.d/fido
service fido
{
        socket_type     = stream
        protocol        = tcp
        wait            = no
        user            = mbse
        instances       = 10
        server          = /opt/mbse/bin/mbcico
        server_args     = -t ifc
}
