# supervisord.d/bbs.ini
[program:mbtask]
command=%(ENV_MBSE_ROOT)s/bin/mbtask -nd
user=mbse
