# supervisord.d/web-server.ini
[program:darkhttpd]
command=/usr/sbin/darkhttpd /srv/mbse/www/htdocs
