# supervisord.d/xinetd.ini
[program:xinetd]
command=/usr/sbin/xinetd -dontfork
